# system-integration/systemd/ss-susi-linux.service.in
[Unit]
Description=kasper Linux
Wants=network-online.target
After=network-online.target ss-kasper-server.service

[Service]
Type=idle
ExecStart=@kasperDIR@/bin/kasper-linux start
Restart=on-failure

[Install]
WantedBy=multi-user.target

